# etc/systemd/system/tg-bot.service
[Unit]
Description=Telegram VPN Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/tg_bot
ExecStart=/opt/tg_bot/venv/bin/python /opt/tg_bot/bot.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
