# etc/systemd/system/tier1-alert.timer
[Unit]
Description=Run Tier-1 Telegram CRIT alert every minute

[Timer]
OnBootSec=45
OnUnitActiveSec=60
AccuracySec=5
Persistent=true

[Install]
WantedBy=timers.target
